# supervisord.conf
; Supervisord configuration for OPB Bot v2.43
; ─────────────────────────────────────────────────────────────────────────────
; Manages two processes:
;   opb_bot        — the main trading bot (index_trader.py)
;   opb_dashboard  — optional FastAPI web dashboard (web_dashboard.py)
;
; The dashboard process is set to autostart=false by default — it is started
; by the bot itself when web_dashboard_enabled=true in config.
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
user=opb

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock


; ── Main bot ──────────────────────────────────────────────────────────────────
[program:opb_bot]
command=python index_app/index_trader.py --paper
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10

; Restart with exponential back-off cap at 300 s
stopwaitsecs=30
killasgroup=true
stopasgroup=true

stdout_logfile=/data/logs/opb_bot.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/data/logs/opb_bot_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

environment=PYTHONUNBUFFERED="1",TZ="Asia/Kolkata"
priority=10


; ── Schema generation (one-shot, runs once at startup) ───────────────────────
[program:opb_schema_gen]
command=python scripts/generate_config_schemas.py
directory=/app
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/data/logs/schema_gen.log
stderr_logfile=/data/logs/schema_gen.log
priority=5


; ── Web dashboard (managed externally by maybe_start_dashboard()) ─────────────
; Set autostart=false — the bot calls core.web_dashboard.maybe_start_dashboard()
; which launches a daemon thread internally.  This entry is kept as a fallback
; for operators who prefer supervisor to manage it directly.
[program:opb_dashboard]
command=python -c "
import json, pathlib, core.web_dashboard as d
cfg = json.loads(pathlib.Path('config.json').read_text()) if pathlib.Path('config.json').exists() else {}
app = d.create_app(cfg)
d.serve(app, host=cfg.get('web_dashboard_host','0.0.0.0'), port=int(cfg.get('web_dashboard_port',8765)))
import time
while True: time.sleep(3600)
"
directory=/app
autostart=false
autorestart=true
startretries=3
startsecs=5

stdout_logfile=/data/logs/dashboard.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/dashboard_err.log
stderr_logfile_maxbytes=5MB

environment=PYTHONUNBUFFERED="1",TZ="Asia/Kolkata"
priority=20
